# roles/logrotate/templates/logrotate.timer.j2
[Unit]
Description=Daily rotation of log files

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
